# services/security-capture-cleanup.timer
[Unit]
Description=Run security-core packet capture cleanup daily

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
RandomizedDelaySec=10m
Unit=security-capture-cleanup.service

[Install]
WantedBy=timers.target
