# services/security-inventory-enrich.timer
[Unit]
Description=Run security-core inventory enrichment every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=security-inventory-enrich.service

[Install]
WantedBy=timers.target
