# services/security-upnp-monitor.service
[Unit]
Description=security-core UPnP/NAT-PMP monitor
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=securitycore
Group=securitycore
EnvironmentFile=/etc/security-core/security-core.env
Environment=PYTHONPATH=/opt/security-core/app/worker:/opt/security-core/app
WorkingDirectory=/opt/security-core/app/worker
ExecStart=/opt/security-core/venv/bin/python /opt/security-core/app/worker/upnp_monitor.py
