# services/security-report-weekly.service
[Unit]
Description=security-core Phase 7 weekly report generation
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=oneshot
User=securitycore
Group=securitycore
EnvironmentFile=/etc/security-core/security-core.env
Environment=PYTHONPATH=/opt/security-core/app/worker:/opt/security-core/app
WorkingDirectory=/opt/security-core/app/worker
ExecStart=/opt/security-core/venv/bin/python /opt/security-core/app/worker/report_engine.py generate --period-days 7 --format html --report-type weekly --actor systemd-weekly
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
