# services/security-opnsense-flow-ingest.timer
[Unit]
Description=Run security-core OPNsense traffic flow ingest every minute

[Timer]
OnBootSec=45s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true
Unit=security-opnsense-flow-ingest.service

[Install]
WantedBy=timers.target
