# services/security-adguard-ingest.timer
[Unit]
Description=Run security-core AdGuard ingest every minute

[Timer]
OnBootSec=100s
OnUnitActiveSec=60s
Persistent=true
RandomizedDelaySec=10s
Unit=security-adguard-ingest.service

[Install]
WantedBy=timers.target
