# services/security-policy-enforcer.timer
[Unit]
Description=Run security-core policy enforcer every minute

[Timer]
OnBootSec=45s
OnUnitActiveSec=60s
Persistent=true
Unit=security-policy-enforcer.service

[Install]
WantedBy=timers.target
