# services/security-detection-correlator.timer
[Unit]
Description=Run security-core detection correlator every two minutes

[Timer]
OnBootSec=90s
OnUnitActiveSec=120s
AccuracySec=15s
Persistent=true
Unit=security-detection-correlator.service

[Install]
WantedBy=timers.target
