# services/security-capture-finalize.timer
[Unit]
Description=Run security-core packet capture processor every 10 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=10s
AccuracySec=1s
Persistent=true
Unit=security-capture-finalize.service

[Install]
WantedBy=timers.target
